# conf.list
/etc/sysconfig/keyboard
/etc/sysconfig/i18n
/etc/sysconfig/clock
/etc/selinux/config
/etc/sysconfig/network
/etc/group
/etc/hosts
/etc/nsswitch.conf
/etc/sysconfig/snmpd
/etc/snmp/snmpd.conf
/etc/udev/rules.d/70-persistent-net.rules
/etc/sysconfig/network-scripts/ifcfg-eth0
/etc/sysconfig/network-scripts/route-eth0
/etc/passwd
/etc/sudoers
/etc/fstab
/boot/grub/grub.conf
/etc/sysctl.conf
/etc/modprobe.d/disable_ipv6.conf
/etc/abrt/abrt-action-save-package-data.conf
/etc/abrt/abrt.conf
/etc/libreport/events.d/ccpp_event.conf
/etc/libreport/events.d/abrt_event.conf
/etc/abrt/plugins/CCpp.conf
/etc/init.d/functions
/etc/security/limits.conf
/etc/security/limits.d/90-nproc.conf
/etc/pam.d/login
/etc/pam.d/password-auth-ac
/etc/pam.d/su
/etc/pam.d/sshd
/etc/ntp.conf
/etc/sysconfig/ntpd
/etc/inittab
/etc/rsyslog.conf
/etc/security/access.conf
/etc/ssh/sshd_config
/etc/aliases
/etc/fstab